# module/katashie-core-bot.service
[Unit]
Description=KATASHIE VPN — Core Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/etc/katashie_bot
ExecStart=/usr/bin/python3 /etc/katashie_bot/katashie_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
